# poky/meta-squeezeos/packages/base-files/files/inittab
::sysinit:/etc/init.d/rcS
::askfirst:/bin/sh
tty3::respawn:/sbin/getty tty3 9600 VC vt100
